# deploy/tyca-tool.service
[Unit]
Description=TYCA Teaching Tool Backend
Wants=network-online.target vpn-codemao.service
After=network-online.target vpn-codemao.service

[Service]
Type=simple
User=tyca-tool
Group=tyca-tool
WorkingDirectory=/opt/tyca-tool/app
EnvironmentFile=/opt/tyca-tool/server.env
ExecStartPre=/usr/sbin/ip link show tun0
ExecStart=/usr/bin/python3.11 /opt/tyca-tool/app/server/app.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/tyca-tool/data /opt/tyca-tool/app/server/data /opt/tyca-tool/录题助手v5.19.11

[Install]
WantedBy=multi-user.target
